# deploy/grapheon-upgrade.path
[Unit]
Description=Watch for Grapheon upgrade requests
Documentation=https://github.com/BadgerOps/grapheon

[Path]
PathExists=/data/upgrade-requested
MakeDirectory=yes

[Install]
WantedBy=multi-user.target
