# deploy/grapheon-upgrade.service
[Unit]
Description=Grapheon container upgrade
Documentation=https://github.com/BadgerOps/grapheon
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/grapheon/scripts/grapheon-upgrade.sh
Environment=DATA_DIR=/data
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

# Run as root to manage systemd services and podman
User=root
Group=root

[Install]
WantedBy=multi-user.target
